# cloudflare-ddns.timer
[Unit]
Description=Timer para Cloudflare DDNS Updater

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
